# system/etc/systemd/system/display.service.d/10-dkms-ensure.conf
# Make jefflessard's display.service wait for the DKMS self-heal, so the clock
# daemon only starts once the module is built and the panel device exists.
[Unit]
Requires=tm16xx-dkms.service
After=tm16xx-dkms.service
